# etc/init.d/xbmc
#!/bin/bash
#
#Copyright 2012 CurlyMo <development@xbian.org>
#
#This file is part of XBian - XBMC on the Raspberry Pi.
#
#XBian is free software: you can redistribute it and/or modify it under the
#terms of the GNU General Public License as published by the Free Software
#Foundation, either version 3 of the License, or (at your option) any later
#version.
#
#XBian is distributed in the hope that it will be useful, but WITHOUT ANY
#WARRANTY; without even the implied warranty of MERCHANTABILITY or FITNESS
#FOR A PARTICULAR PURPOSE. See the GNU General Public License for more
#details.
#
#You should have received a copy of the GNU General Public License along
#with XBian. If not, see <http://www.gnu.org/licenses/>
#
### BEGIN INIT INFO
# Provides:          xbmc
# Required-Start:    bootlocal
# Required-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start and stop XBMC through init
# Description:
#
### END INIT INFO

DESC=XBMC
NAME=xbmc.bin
DAEMON=/usr/local/lib/xbmc/$NAME
DAEMON_ARGS="--standalone"
PIDFILE=/var/run/xbmc.pid
SCRIPTNAME=/etc/init.d/xbmc

. /lib/init/vars.sh
. /lib/lsb/init-functions

# Exit if the package is not installed
[ -x $DAEMON ] || exit 0

do_start()
{
	if [ $(ps -A | grep xbmc.bin | wc -l) -eq 0 ]; then
		/usr/bin/black;
		start-stop-daemon -c xbian -u xbian --start --quiet --pidfile $PIDFILE --exec $DAEMON --test || return 1;
		echo $(start-stop-daemon -c xbian -u xbian -m --start --quiet --pidfile $PIDFILE --exec $DAEMON -- $DAEMON_ARGS; RETURN=$?; case $RETURN in 0 ) pkill black; ;; 64 ) sudo halt ;; 66 ) sudo reboot; ;; esac || return 2) &exit
	fi
}

do_stop()
{
	start-stop-daemon --stop --signal 9 --retry 5 --quiet --pidfile $PIDFILE
	pkill black;
}

case "$1" in
  start)
        [ "$VERBOSE" != no ] && log_daemon_msg "Starting $DESC" "$NAME"
        do_start
        case 0 in
                0|1) [ "$VERBOSE" != no ] && log_end_msg 0 ;;
                2) [ "$VERBOSE" != no ] && log_end_msg 1 ;;
        esac
        ;;
  stop)
        [ "$VERBOSE" != no ] && log_daemon_msg "Stopping $DESC" "$NAME"
        do_stop
        case "$?" in
                0|1) [ "$VERBOSE" != no ] && log_end_msg 0 ;;
                2) [ "$VERBOSE" != no ] && log_end_msg 1 ;;
        esac
        ;;
  restart)
        $0 stop
        sleep 1
        $0 start
        ;;
  status)
        status_of_proc "$DAEMON" "$NAME" && exit 0 || exit 0
        ;;
  *)
        echo "Usage: $SCRIPTNAME {start|stop|restart|status}" >&2
        exit 3
        ;;
esac

:
